# client/autostart/agentos-voice.service
# systemd user service (Linux laptop). Install:
#   mkdir -p ~/.config/systemd/user
#   cp client/autostart/agentos-voice.service ~/.config/systemd/user/
#   edit paths/env below, then:
#   systemctl --user daemon-reload && systemctl --user enable --now agentos-voice
[Unit]
Description=AgentOS voice client
After=default.target sound.target

[Service]
Type=simple
WorkingDirectory=%h/AgentOS/client
Environment=AGENT_WS_URL=wss://agent.example.com/ws/voice
Environment=AGENT_WS_TOKEN=CHANGE_ME
ExecStart=%h/AgentOS/.venv/bin/python voice_client.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
